# etc/systemd/birdwatcher.target
[Unit]
Wants=birdwatcher4.service birdwatcher6.service

[Install]
WantedBy=multi-user.target

