# install/upstart/init/birdwatcher.conf

#
# Birdwatcher Startup
#

description "birdwatcher"
author "Matthias Hannig <mha@ecix.net>"


start on runlevel [2345]
stop on runlevel [!2345]


